# intelli-credit/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
loglevel=info

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:backend]
command=node src/index.js
directory=/app/backend
autostart=true
autorestart=true
environment=PORT="3001",GO_SERVICE_URL="http://127.0.0.1:8081",AI_SERVICE_URL="http://127.0.0.1:8000",SHARED_TMP_PATH="/tmp/intelli-credit"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:go-service]
command=/app/go-service/go-service
directory=/app/go-service
autostart=true
autorestart=true
environment=PORT="8081"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:ai-service]
command=uvicorn main:app --host 0.0.0.0 --port 8000
directory=/app/ai-service
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
